# apps/api/supervisord.conf
; supervisord.conf — manages cloudflared + strapi in a single container
;
; Required environment variable:
;   CLOUDFLARE_TUNNEL_TOKEN — token from Cloudflare Zero Trust dashboard
;     (Tunnels → <tunnel> → Configure → Install connector → Docker token)

[supervisord]
nodaemon=true
; Write logs to /dev/null — each program writes directly to stdout/stderr
logfile=/dev/null
logfile_maxbytes=0
; /tmp is world-writable so the strapi user can create the pidfile
pidfile=/tmp/supervisord.pid

; ── Cloudflare Tunnel ─────────────────────────────────────────────────────────
; Establishes an outbound-only tunnel to Cloudflare.
; Proxies external HTTPS traffic to Strapi at localhost:1337.
; Starts first (priority=10) so the tunnel is ready before Strapi boots.
[program:cloudflared]
command=cloudflared tunnel run --token %(ENV_CLOUDFLARE_TUNNEL_TOKEN)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Strapi ────────────────────────────────────────────────────────────────────
; Starts the production Strapi server on localhost:1337.
[program:strapi]
command=pnpm start
directory=/app/apps/api
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
